# ops/github-release-updater/areaforge-updater.timer
[Unit]
Description=Check AreaForge GitHub Releases for updates

[Timer]
OnCalendar=*-*-* 04:20:00
RandomizedDelaySec=1800
Persistent=true

[Install]
WantedBy=timers.target
